# scripts/wazuh-api
#!/bin/sh
# WAZUH API Service
# Author: Wazuh

### BEGIN INIT INFO
# Provides:          wazuh_api
# Required-Start:    $remote_fs $syslog
# Required-Stop:     $remote_fs $syslog
# Should-Start:      $network
# Should-Stop:       $network
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Wazuh API
# Description:       Wazuh API daemon
#
### END INIT INFO

# Instructions:
#  sudo install -m 755 -o root -g root /var/ossec/api/scripts/wazuh-api /etc/init.d/
#  sudo systemctl enable wazuh-api
#  sudo systemctl daemon-reload
#  sudo systemctl restart wazuh-api

. /etc/ossec-init.conf

if [ "X${DIRECTORY}" = "X" ]; then
    DIRECTORY="/var/ossec"
fi

BIN_DIR=$(which nodejs 2> /dev/null)

if [ "X$BIN_DIR" = "X" ]; then
    BIN_DIR=$(which node 2> /dev/null)

    if [ "X$BIN_DIR" = "X" ]; then
        echo "NodeJS binaries not found"
        exit 1
    fi
fi

PID_DIR="/var/run/wazuh-api.pid"
APP_PATH="${DIRECTORY}/api/app.js"
LOG_DIR="${DIRECTORY}/logs/api.log"

start() {
	if [ ! -f $PID_DIR ]; then
		($BIN_DIR $APP_PATH >> $LOG_DIR) & echo $!>$PID_DIR
	fi
}

stop() {
	if [ -f $PID_DIR ]; then
		kill `cat $PID_DIR`
	    rm $PID_DIR
	fi
}


case "$1" in
start)
    start
    ;;
stop)
    stop
    ;;
restart)
    stop
    start
    ;;
status)

    ;;
*)
    echo "*** Usage: $0 {start|stop|restart|status}"
    exit 1
esac
